# scripts/moode-extmgr.service
[Unit]
Description=moOde ext-mgr control plane
After=network.target
Wants=moode-extmgr-watchdog.service

[Service]
Type=simple
User=moode-extmgrusr
Group=moode-extmgr
ExecStart=/usr/local/bin/moode-extmgr-service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
