# inspiration/Test extensions/ronnie-pickering-extension-1.0.0/RonniePickeringExtension/packages/services/ronnie-pickering-helper.service
[Unit]
Description=Ronnie Pickering Extension – helper service (test package)
After=network.target moode-extmgr.service
Requires=moode-extmgr.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=moode-extmgrusr
ExecStart=/bin/true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
